# crontab
0 */1 * * * docker exec kata_app_2 rails runner /app/script/update_hot_posts.rb
